# kmod/evo4-load-config.service
[Unit]
Description=Load Audient EVO4 configuration on device connection
ConditionPathExists=/dev/evo4

[Service]
Type=oneshot
ExecStart=%h/.local/bin/evoctl --device evo4 load
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
