# kmod/evo8-load-config.service
[Unit]
Description=Load Audient EVO8 configuration on device connection
ConditionPathExists=/dev/evo8

[Service]
Type=oneshot
ExecStart=%h/.local/bin/evoctl --device evo8 load
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
